# rootfs/etc/s6-overlay/s6-rc.d/upsd/run
#!/command/with-contenv sh

sleep 3
echo "+2.正在启动 upsd 服务端...启动命令：\"upsd $UPSD_OPTS\""
exec upsd $UPSD_OPTS