# rootfs/etc/s6-overlay/s6-rc.d/upsd/run
#!/command/with-contenv sh

echo "+正在启动 upsd 服务端..."
# 启动服务端 调试+前台模式
exec upsd -D