# rootfs/etc/s6-overlay/s6-rc.d/darkhttpd/finish
#!/command/with-contenv sh

if [ "$WEB" == "true" ]; then
  echo "正在停止 darkhttpd 服务..."
  pkill darkhttpd
else
  s6-svc -Od .
  exit 0
fi