# rootfs/etc/s6-overlay/s6-rc.d/caddy/finish
#!/command/with-contenv sh

echo "-正在停止 Caddy 服务..."