# rootfs/etc/s6-overlay/s6-rc.d/init/up
#!/command/execlineb

foreground {
    /etc/s6-overlay/s6-rc.d/init/08-init.sh
}
/etc/s6-overlay/s6-rc.d/init/00-banner.sh