# rootfs/etc/s6-overlay/s6-rc.d/nut/finish
#!/command/with-contenv sh

echo "+正在清理PID文件..."

rm -rf /var/run/nut/*

# 停止容器
/run/s6/basedir/bin/halt