# rootfs/etc/s6-overlay/s6-rc.d/upsmon/run
#!/command/with-contenv sh

sleep 1
echo "+3.正在启动 upsmon 监控服务..."
# 启动监控服务 调试+前台模式
exec upsmon -D