# rootfs/etc/s6-overlay/s6-rc.d/upsmon/run
#!/command/with-contenv sh

sleep 1
echo "+3.正在启动 upsmon 监控服务..."
echo "启动命令：\"upsmon $UPSMON_OPTS\""
exec upsmon $UPSMON_OPTS