# rootfs/etc/s6-overlay/s6-rc.d/geoipupdate/finish
#!/command/with-contenv sh

echo "-正在停止 GeoIP Update 服务..."