# rootfs/etc/s6-overlay/s6-rc.d/vlmcsd/run
#!/command/with-contenv sh

#启动vlmcsd
exec s6-setuidgid kms vlmcsd -i /vlmcsd/vlmcsd.ini -D