# rootfs/etc/s6-overlay/s6-rc.d/vlmcsd/run
#!/bin/sh

#启动vlmcsd
exec s6-setuidgid kms vlmcsd -i /etc/vlmcsd.ini -D