# rootfs/etc/s6-overlay/s6-rc.d/vlmcsd/run
#!/command/with-contenv sh

echo "正在启动 vlmcsd 服务..."
exec s6-setuidgid kms vlmcsd -i /vlmcsd/vlmcsd.ini -D