# rootfs/etc/s6-overlay/s6-rc.d/upsdrvctl/finish
#!/command/with-contenv sh

echo "+正在清理PID文件..."
rm -rf /var/run/nut/*