# rootfs/etc/s6-overlay/s6-rc.d/upsdrvctl/run
#!/command/with-contenv sh

echo "+1.正在启动 upsdrvctl 驱动服务..."
# 启动驱动 调试+前台模式
exec upsdrvctl -D -FF start