# rootfs/etc/s6-overlay/s6-rc.d/upsdrvctl/run
#!/command/with-contenv sh

echo "+1.正在启动 upsdrvctl 驱动服务...启动命令：\"upsdrvctl $UPSDRVCTL_OPTS\""
exec upsdrvctl $UPSDRVCTL_OPTS start